# crontab
SHELL=/bin/bash
* * * * * /home/ubuntu/whisper-server-script/run.sh 2>>/home/ubuntu/whisperx.log 1>>/home/ubuntu/whisperx.log
